# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
autorestart=true

[program:powershell-script]
command=pwsh /app/script.ps1
autorestart=true
